# systemd/rpitank.service
[Unit]
Description=RPiTank control server
After=mjpg_streamer.service

[Service]
ExecStart=/usr/bin/rpitank

[Install]
WantedBy=multi-user.target
